# vmu-packer.cron
30 0 * * 7 root /usr/bin/vmu-rebuild-all
